# post-install/dotfiles/config/dotfiles-pull.service
[Unit]
Description=Auto-pull dotfiles from git repository
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/git-autopull --repo dotfiles
Restart=on-failure
Environment=DISPLAY=:0

[Install]
WantedBy=default.target